# .config/autostart/kdeconnectd.desktop
[Desktop Entry]
Type=Application
Name=KDE Connect Daemon
Exec=/usr/lib/kdeconnectd
X-GNOME-Autostart-enabled=true
